# apps/backend/cloudflared-gui-backend.service
[Unit]
Description=Cloudflared GUI Backend Service
After=network.target cloudflared.service
Wants=cloudflared.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cloudflared-gui/backend
ExecStart=/opt/cloudflared-gui/backend/server
Restart=on-failure
RestartSec=5s

# Environment variables (optional - defaults shown)
# Uncomment and modify as needed:
# Environment="HOST=127.0.0.1"
# Environment="PORT=8080"
# Environment="CLOUDFLARED_SERVICE_NAME=cloudflared.service"
# Environment="CLOUDFLARED_CONFIG_PATH=/etc/cloudflared/config.yml"

# Security settings
NoNewPrivileges=false
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared-gui-backend

[Install]
WantedBy=multi-user.target

